# deploy/launchd/com.zeus.loop-tick.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Zeus 24/7 improvement loop v3 — single tick trigger.
  Authority: docs/operations/current/plans/allday_improvement_loop_v3_codex_2026-07-09.md.

  TEMPLATE ONLY — not installed by any script in this repo. No secrets.

  Fires loop/tick.sh hourly at :17 — the FINEST granularity, not the
  effective cadence. loop/INTERVAL (a number of hours, operator-edited,
  default 1) gates each firing inside tick.sh; retuning the loop to 2h/6h/
  anything never touches this plist. StartCalendarInterval (not
  StartInterval) so it neither drifts nor double-fires on wake-from-sleep.

  PATH must carry both the repo venv and the codex CLI
  (~/.npm-global/bin). The engine is `codex exec` under an OS sandbox —
  see loop/tick.sh.

  INSTALL (operator, after scripts/ops/loop_enable.sh preflight passes):
    1. cp deploy/launchd/com.zeus.loop-tick.plist ~/Library/LaunchAgents/
    2. launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.zeus.loop-tick.plist

  DISABLE: `touch <repo>/loop/HALT` (soft — daemon stays loaded, next fire
  no-ops) or `launchctl bootout gui/$(id -u)/com.zeus.loop-tick` (hard).
-->
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/leofitz/zeus/.venv/bin:/Users/leofitz/.npm-global/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>PYTHONPATH</key>
		<string>/Users/leofitz/zeus</string>
		<key>ZEUS_PRIMARY_ROOT</key>
		<string>/Users/leofitz/zeus</string>
	</dict>
	<key>Label</key>
	<string>com.zeus.loop-tick</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/leofitz/zeus/loop/tick.sh</string>
	</array>
	<key>StandardErrorPath</key>
	<string>/Users/leofitz/zeus/loop/logs/launchd-tick.err</string>
	<key>StandardOutPath</key>
	<string>/Users/leofitz/zeus/loop/logs/launchd-tick.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Minute</key>
		<integer>17</integer>
	</dict>
	<key>ThrottleInterval</key>
	<integer>60</integer>
	<key>WorkingDirectory</key>
	<string>/Users/leofitz/zeus</string>
</dict>
</plist>
