# deploy/launchd/com.zeus.post-trade-capital.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Created: 2026-06-08
  Authority basis: docs/architecture/system_decomposition_plan.md §4.3 / §6 (P4 row) / §8 Step 2.

  ARTIFACT ONLY. This launchd job description is committed for the operator to install
  when the cutover is performed; the process-topology refactor that produced it does NOT
  load, kickstart, or enable any service. It mirrors the existing com.zeus.data-ingest /
  com.zeus.substrate-observer pattern (KeepAlive + RunAtLoad + ThrottleInterval,
  ProgramArguments -> `python -m <daemon module>`).

  This is the NEW P4 post-trade-capital program: it runs the lifted POST_TRADE capital
  follow-up cycles (chain-sync READ phase + harvester resolver + redeem x2 + wrap x3) that
  drive the settlement_commands / wrap_unwrap_commands state machines. The order daemon
  (com.zeus.live-trading) no longer schedules these jobs; it keeps ONLY the exit-monitoring
  / exit-SUBMIT phase (which posts real sell orders) — §8 Step 2.

  INSTALL (operator, NOT performed by the refactor):
    cp deploy/launchd/com.zeus.post-trade-capital.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.zeus.post-trade-capital.plist
-->
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>/Users/leofitz/zeus</string>
		<key>PYTHONSAFEPATH</key>
		<string>1</string>
		<key>POLYMARKET_CLOB_V2_SIGNATURE_TYPE</key>
		<string>2</string>
		<key>ZEUS_CLOB_FACT_CACHE_TTL_SECONDS</key>
		<string>0</string>
		<key>ZEUS_PRIMARY_ROOT</key>
		<string>/Users/leofitz/zeus</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.zeus.post-trade-capital</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/leofitz/zeus/.venv/bin/python</string>
		<string>-m</string>
		<string>src.ingest.post_trade_capital_daemon</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/leofitz/zeus/logs/zeus-post-trade-capital.err</string>
	<key>StandardOutPath</key>
	<string>/Users/leofitz/zeus/logs/zeus-post-trade-capital.log</string>
	<key>ThrottleInterval</key>
	<integer>30</integer>
	<key>WorkingDirectory</key>
	<string>/Users/leofitz/zeus</string>
</dict>
</plist>
